# systemd/whisper-stt.service.d/gpu-fleet-lease.conf
# gpu-fleet onboarding for whisper-stt (RFC 0006 / GPUFLE-1): while STT is hot,
# this box's llama slot row is exclusively leased, so `pick`'s derived capacity
# drops it below routable and fleet consumers SKIP the shared 3090 instead of
# colliding with whisper's ~0.95 GiB (the old CUDA-OOM class).
#
# Install: /etc/systemd/system/whisper-stt.service.d/gpu-fleet-lease.conf
# then: sudo systemctl daemon-reload && sudo systemctl restart whisper-stt
# Renewal lives in the whisper-stt-lease-renew.service companion (same dir).

[Unit]
# acquire exits 75 while another consumer's lease is live (the scheduling skip);
# Restart=on-failure then retries every RestartSec until the lease drains. A
# fleet shard can hold a slot for many minutes, so never rate-limit the retries.
StartLimitIntervalSec=0

[Service]
# Claim the slot BEFORE whisper-server maps VRAM; release (fenced) after it
# stops. Runs as the unit's User= (halbritt), same as the heartbeat writer.
ExecStartPre=/home/halbritt/git/gpu-fleet/bin/whisper-stt-lease acquire
ExecStopPost=/home/halbritt/git/gpu-fleet/bin/whisper-stt-lease release
